# settings/cron.d/mm5d
# /etc/cron.d/mm5d: crontab entries for MM5D

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 1 * * *			root	mm5d-updatestartpage >/dev/null 2>&1
0 2 * * *			root	mm5d-maintainlog >/dev/null 2>&1
